# etc/profile.d/lang.csh
version https://git-lfs.github.com/spec/v1
oid sha256:4dc37ad18435ec9208213c2005fc2602bf29ad58720d04012fa004554843afc7
size 866
