# etc/profile.d/tzset.csh
version https://git-lfs.github.com/spec/v1
oid sha256:05f2eff08911b3b7896f563aa7700e00d6332fa996d873de696d0cb7380b53c7
size 1037
