# etc/defaults/etc/profile.d/lang.sh
version https://git-lfs.github.com/spec/v1
oid sha256:24085cd235392724b2c44163f595044319dab68d01d1b4a69c99e7d16665828f
size 857
